# docker/cron/crontab
0 0 * * * docker exec nginx /scripts/logrotate.sh > /dev/stdout
